# roles/todo1/files/todo1.service
[Unit]
Description = Todo Service
[Service]
User=root
Environment=REDIS_HOST=172.31.9.149
Environment=REDIS_PORT=6379
Environment=USERS_API_ADDRESS=http://172.31.12.70:8080
WorkingDirectory=/root/todo1
ExecStart=npm start
SyslogIdentifier=todo1

[Install]
WantedBy=multi-user.target
